# supervisor.conf
[program:banner_websocket]
command=/usr/bin/node /var/www/banner_system/ws.njs
directory=/var/www/banner_system
restart=always
user=root
